# roles/traefik_api/templates/traefik.service.j2
[Unit]
Description=Traefik Reverse Proxy
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User={{ traefik_user }}
Group={{ traefik_group }}
EnvironmentFile=-{{ traefik_config_dir }}/traefik.env
{% set _cli_flags = traefik_cli_flags | default([]) %}
{% set _cli_suffix = ' ' + (_cli_flags | join(' ')) if _cli_flags else '' %}
ExecStart={{ traefik_install_dir }}/traefik --configFile={{ traefik_config_dir }}/traefik.yml{{ _cli_suffix }}
Restart=always
RestartSec=5s
LimitNOFILE=65536
AmbientCapabilities=CAP_NET_BIND_SERVICE
CapabilityBoundingSet=CAP_NET_BIND_SERVICE
NoNewPrivileges=true
ProtectSystem=full
ProtectHome=true
PrivateTmp=true
ReadWritePaths={{ traefik_data_dir }} {{ traefik_log_dir }}
WorkingDirectory={{ traefik_config_dir }}
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
