# roles/nextcloud/templates/nextcloud-cron.service.j2
[Unit]
Description=Nextcloud cron.php job

[Service]
User={{ nc_web_user }}
Group={{ nc_web_group }}
ExecStart=/usr/bin/php -f {{ nc_install_dir }}/cron.php

[Install]
WantedBy=multi-user.target
