# roles/keycloak/templates/keycloak.service.j2
[Unit]
Description=Keycloak Identity Provider
After=network.target
Wants=network-online.target

[Service]
Type=simple
User={{ keycloak_runtime_user }}
Group={{ keycloak_runtime_group }}
EnvironmentFile=-/etc/keycloak/keycloak.env
WorkingDirectory={{ keycloak_install_dir }}
ExecStart={{ keycloak_install_dir }}/bin/kc.sh start --config-file={{ keycloak_install_dir }}/conf/keycloak.conf
Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=120
LimitNOFILE=65536
NoNewPrivileges=true
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
